# muninndb/rootfs/etc/services.d/muninndb/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: MuninnDB
# Take down the S6 supervision tree when MuninnDB fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "MuninnDB crashed with exit code ${1}. Respawning..."
fi
