# huly/rootfs/etc/services.d/huly/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Huly
# Stops the Huly stack cleanly on shutdown or crash
# ==============================================================================
bashio::log.info "Stopping Huly stack..."
cd /data/huly 2>/dev/null && docker compose down 2>/dev/null

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Huly crashed with exit code ${1}. Respawning..."
fi
